# php-fpm.conf
[program:php-fpm]
command = /usr/sbin/php-fpm7.2 -c /etc/php/7.2/fpm --nodaemonize
autostart=true
autorestart=true
